# root/etc/crontabs/root
*/2 * * * * /usr/bin/awg-healthcheck
1-59/2 * * * * /usr/bin/singbox-healthcheck
